# rpm.cron
#!/bin/sh
umask 022
export LC_ALL=C

rpm -qa --qf '%{name}-%{version}-%{release}.%{arch}.rpm\n' 2>&1 \
	| sort > /var/log/rpmpkgs
