# rpm.cron
#!/bin/sh

. /etc/sysconfig/rpm

if [ ${RPM_CRON_RPMQA_CACHE:-no} != "yes" ]; then
   exit 0
fi

umask 022
export LC_ALL=C

rpm -qa --qf '%{name}-%{version}-%{release}.%{arch}.rpm\n' 2>&1 \
	| sort > /var/log/rpmpkgs
